# go/platform/deployment/systemd/iron-signal-foundation-api.service
[Unit]
Description=Iron Signal Platform Foundation API
Documentation=https://github.com/Iron-Signal-Systems/iron-signal-platform
Wants=network-online.target
After=network-online.target
StartLimitIntervalSec=60s
StartLimitBurst=3

[Service]
Type=notify
NotifyAccess=main
User=issp-foundation-api
Group=issp-foundation-api
ExecStart=/usr/lib/iron-signal-platform/foundation-api
LoadCredentialEncrypted=database-url:/etc/iron-signal-platform/credentials/foundation-api.database-url.cred
Environment=ISSP_ADMIN_LISTEN_ADDRESS=127.0.0.1:18081
Environment=ISSP_DATABASE_DSN_FILE=%d/database-url
Environment=ISSP_DATABASE_CONNECT_TIMEOUT=5s
Environment=ISSP_DATABASE_MAX_CONNECTIONS=4
Environment=ISSP_DATABASE_MIN_CONNECTIONS=0
Environment=ISSP_DATABASE_MAX_CONNECTION_LIFETIME=30m
Environment=ISSP_DATABASE_MAX_CONNECTION_IDLE_TIME=5m
Environment=ISSP_DATABASE_HEALTH_CHECK_PERIOD=30s
Environment=ISSP_STARTUP_TIMEOUT=15s
Environment=ISSP_SHUTDOWN_TIMEOUT=10s
Restart=on-failure
RestartSec=5s
RestartPreventExitStatus=78
TimeoutStartSec=30s
TimeoutStopSec=20s
WatchdogSec=30s
KillMode=control-group
KillSignal=SIGTERM
FinalKillSignal=SIGKILL
SendSIGKILL=yes
UMask=0077
NoNewPrivileges=yes
CapabilityBoundingSet=
AmbientCapabilities=
PrivateTmp=yes
PrivateDevices=yes
PrivateMounts=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
ProtectProc=invisible
ProcSubset=pid
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RemoveIPC=yes
KeyringMode=private
InaccessiblePaths=/etc/iron-signal-platform/credentials
LimitNOFILE=1024
TasksMax=64
MemoryHigh=192M
MemoryMax=256M
StandardOutput=journal
StandardError=journal
SyslogIdentifier=iron-signal-foundation-api
LogRateLimitIntervalSec=30s
LogRateLimitBurst=200

[Install]
WantedBy=multi-user.target
